# deploy/gcp/loom-backup.timer
[Unit]
Description=Nightly loom SQLite backup

[Timer]
OnCalendar=*-*-* 03:00:00 UTC
RandomizedDelaySec=45m
Persistent=true
Unit=loom-backup.service

[Install]
WantedBy=timers.target

